# frontend/infrastructure_files/gcs_bucket/clusters/ansible_setup/roles/c2_daemon/files/ghpcfe_c2.service
[Unit]
Description=GHPCFE Command and Control Daemon
After=network.target

[Service]
Type=simple
PIDFile=/run/ghpcefe.pid
ExecStart=/usr/local/sbin/ghpcfe_c2daemon.py
ExecStop=/bin/kill -s INT $MAINPID
RestartForceExitStatus=123

[Install]
WantedBy=multi-user.target
