# community/modules/scripts/vdi-setup/roles/lock_manager/templates/vdi-monitor.service.j2
[Unit]
Description=VDI Configuration Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/vdi-monitor.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
