# community/modules/scheduler/schedmd-slurm-gcp-v6-controller/modules/slurm_files/scripts/slurmhcd.service
[Unit]
Description=slurmhcd (A Slurm Healthcheck Daemon)

[Service]
User=root
ExecStart=/usr/bin/python3 /slurm/scripts/slurmhcd.py
Restart=always

[Install]
WantedBy=multi-user.target
